# scripts/supervisor/rococo.conf
[program:pichiu-collator-1]
command=/data/kylin-collator --collator --name pichiu-collator-1 --force-authoring --parachain-id 2102 --chain pichiu-rococo-parachain-2102-raw.json --base-path  /data/para_db/pichiu_collator_1 --port 40333 --ws-port 8844 --rpc-cors all -- --execution wasm --chain ./rococo.json --port 30343 --ws-port 9977
directory=/data
autostart=true                ; start at supervisord start (default: true)
autorestart=true        ; whether/when to restart (default: unexpected)
startsecs=1                   ; number of secs prog must stay running (def. 1)
startretries=3                ; max # of serial start failures (default 3)
user=ubuntu                     ; setuid to this UNIX account to run the program
redirect_stderr=true          ; redirect proc stderr to stdout (default false)
stdout_logfile=/data/logs/collator-01.log
stopasgroup=true
nodaemon=true

[program:pichiu-collator-2]
command=/data/kylin-collator --collator --name pichiu-collator-2 --force-authoring --parachain-id 2102 --chain pichiu-rococo-parachain-2102-raw.json --base-path  /data/para_db/pichiu_collator_2 --port 40334 --ws-port 8845 --rpc-cors all -- --execution wasm --chain ./rococo.json --port 30344 --ws-port 9978
directory=/data
autostart=true                ; start at supervisord start (default: true)
autorestart=true        ; whether/when to restart (default: unexpected)
startsecs=1                   ; number of secs prog must stay running (def. 1)
startretries=3                ; max # of serial start failures (default 3)
user=ubuntu                     ; setuid to this UNIX account to run the program
redirect_stderr=true          ; redirect proc stderr to stdout (default false)
stdout_logfile=/data/logs/collator-02.log
stopasgroup=true
nodaemon=true
